# scripts/system/reverse-ssh-daemon.service
[Unit]
Description=ViewTouch Reverse SSH Tunnel Daemon
After=network.target
Wants=network.target

[Service]
Type=simple
User=viewtouch
Group=viewtouch
ExecStart=/usr/viewtouch/bin/reverse_ssh_daemon
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=reverse-ssh-daemon

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectHome=yes
ReadWritePaths=/var/log/viewtouch /var/run/viewtouch /usr/viewtouch/ssh
ProtectSystem=strict
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes

# Resource limits
LimitNOFILE=1024
MemoryLimit=50M

[Install]
WantedBy=multi-user.target
